# docker/root/etc/s6-overlay/s6-rc.d/fix-permissions/up
#!/command/execlineb -P

with-contenv

if {
  bash -c "
    source /usr/local/share/switch-user.sh
    patchUser user "$USER_UID" "$USER_GID"
  "
}
