# examples/chromium/root/etc/s6-overlay/s6-rc.d/chromium/run
#!/command/execlineb -P

with-contenv

redirfd -w 2 /dev/null

s6-env HOME=/home/user

s6-setuidgid user chromium --no-sandbox --start-maximized https://github.com/kyzima-spb/docker-gui/
