# examples/chromium/root/etc/s6-overlay/s6-rc.d/chromium/run
#!/command/execlineb -P

with-contenv

redirfd -w 2 /dev/null

chromium --no-sandbox --start-maximized