# docker/root/etc/s6-overlay/s6-rc.d/openbox-config/up
#!/command/execlineb -P

with-contenv

if {
    xmlstarlet ed -L \
        -N o="http://openbox.org/3.4/rc" \
        -u /o:openbox_config/o:theme/o:name -v Numix \
        -u /o:openbox_config/o:desktops/o:number -v 1 \
            /etc/xdg/openbox/rc.xml
}

foreground {
    redirfd -a 1 /etc/gtk-3.0/settings.ini s6-echo "\n[Settings]\ngtk-theme-name=Numix"
}

redirfd -a 1 /etc/gtk-2.0/gtkrc s6-echo "\ngtk-theme-name=\"Numix\""
