# docker/root/etc/s6-overlay/s6-rc.d/openbox-config/up
#!/command/execlineb -P

with-contenv

importas -u home_dir HOME

multisubstitute
{
    define openbox_config_file ${home_dir}/.config/openbox/rc.xml
    define gtk3_settings_file ${home_dir}/.config/gtk-3.0/settings.ini
    define gtk2_settings_file ${home_dir}/.gtkrc-2.0
}

if {
    mkdir -p ${home_dir}/.config/gtk-3.0 ${home_dir}/.config/openbox
}

if {
    redirfd -w 1 ${openbox_config_file} xmlstarlet ed \
        -N o="http://openbox.org/3.4/rc" \
        -u /o:openbox_config/o:theme/o:name -v Numix \
        -u /o:openbox_config/o:desktops/o:number -v 1 \
            /etc/xdg/openbox/rc.xml
}

foreground {
    redirfd -w 1 ${gtk3_settings_file} s6-echo "[Settings]\ngtk-theme-name=Numix"
}

redirfd -w 1 ${gtk2_settings_file} s6-echo "gtk-theme-name=\"Numix\""
