# examples/kodi/root/etc/s6-overlay/s6-rc.d/kodi/run
#!/command/execlineb -P

with-contenv

backtick -E HOME { homeof user }
s6-env HOME="$HOME"

redirfd -w 2 /dev/null

s6-setuidgid user kodi -fs --standalone
