# docker/root/etc/s6-overlay/s6-rc.d/x11vnc-password/up
#!/command/execlineb -P

with-contenv

backtick -E HOME { homeof user }
backtick -E password { read-secret VNC_SERVER_PASSWORD "" }
define password_file "${HOME}/.vncpswd"

ifthenelse { test -z "$password" }
{
    rm -f "$password_file"
}
{
    execline-umask 600
    redirfd -w 2 /dev/null
    if { s6-setuidgid user /usr/bin/x11vnc -storepasswd "$password" "$password_file" }
}
