# examples/kodi/root/etc/s6-overlay/s6-rc.d/pulse-config/up
#!/command/execlineb -P

with-contenv

if {
    redirfd -r 0 /etc/pulse/client.conf.tmpl
    redirfd -w 1 /etc/pulse/client.conf
    envsubst
}